# scripts/estate-tunnel.service
[Unit]
Description=Cloudflare named tunnel sahiix-estate → :3001
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.cloudflared/sahiix-estate.env
ExecStart=%h/.local/bin/cloudflared tunnel --no-autoupdate run
Restart=always
RestartSec=5
StandardOutput=append:%h/estate-tunnel.log
StandardError=append:%h/estate-tunnel.log

[Install]
WantedBy=default.target
